# templates/dcos-clear-logs.service.j2
[Unit]
Description=Clear DC/OS Logs Service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart dcos-journalctl-filebeat.service
