# templates/dcos-clear-logs.timer.j2
[Unit]
Description=Clear DC/OS Logs Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timer.target
